# systemd/user/backup.service
[Unit]
Description=Restic backup

[Service]
ExecStart=/usr/bin/fish -c backup

[Install]
WantedBy=default.target
